# scripts/simrad-timer.service
[Unit]
Description=Simrad Race Timer NMEA 2000 → Signal K bridge
After=signalk.service can-interface.service
Wants=signalk.service can-interface.service

[Service]
Type=simple
User=helmlog
Group=mark

# Inherit project dir from environment; override in /etc/systemd/system/simrad-timer.service.d/
WorkingDirectory=/home/mark/helmlog

# Read SK token and CAN channel from the same .env as helmlog
EnvironmentFile=-/home/mark/helmlog/.env

ExecStart=/home/mark/.local/bin/uv run --no-sync --project /home/mark/helmlog \
    python scripts/simrad_timer_sk.py \
    --channel ${CAN_INTERFACE:-can0} \
    --signalk ${SK_URL:-http://localhost:3000} \
    --token ${SK_TOKEN:-} \
    --log-level INFO

Restart=on-failure
RestartSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=simrad-timer

[Install]
WantedBy=multi-user.target
